# src/intel/igs/intel-gpu-service.service.in
[Unit]
Description=Intel GPU Service

[Service]
Type=dbus
BusName=com.intel.GPU
ExecStart=@install_libexecdir@/intel-gpu-service
